# supervisor.conf
[supervisord]
nodaemon=true

[program:frontend]
directory= /cicd-accelerator-frontend
command= /usr/bin/npm start


[program:backend]
directory= /cicd-accelerator-backend
command= /usr/bin/npm start